# deploy/systemd/aegis-gpu.service
[Unit]
Description=Aegis-GPU AIOps Backend API
After=network.target
Wants=network.target

[Service]
Type=notify
NotifyAccess=all
User=aegis
Group=aegis
WorkingDirectory=/opt/aegis-gpu
EnvironmentFile=/etc/aegis-gpu/aegis.env
Environment=PYTHONUNBUFFERED=1
Environment=AEGIS_BIND_HOST=127.0.0.1
Environment=AEGIS_BIND_PORT=8000
Environment=AEGIS_WORKERS=4
Environment=AEGIS_READY_TIMEOUT=45
UMask=0027

ExecStart=/usr/bin/python3 /opt/aegis-gpu/aegis_service_runner.py
Restart=always
RestartSec=5
TimeoutStartSec=60
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aegis-gpu

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/aegis-gpu /var/lib/aegis-gpu
ReadOnlyPaths=/opt/aegis-gpu /etc/aegis-gpu

[Install]
WantedBy=multi-user.target
